# aquatiq-webhook.service
[Unit]
Description=Aquatiq Root Container Webhook Deployment Listener
Documentation=https://github.com/I-Dacosta/Aquatiq
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aquatiq
Environment="GITHUB_WEBHOOK_SECRET=%i"
Environment="WEBHOOK_PORT=9000"
ExecStart=/opt/aquatiq/webhook-deploy.sh server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/aquatiq /var/log

[Install]
WantedBy=multi-user.target
